# pkgbuild/pkg/biglinux-improve-compatibility/usr/lib/systemd/system/glibc2-41-fix.path
[Unit]
Description=Fix apps broken with this bug: https://sourceware.org/bugzilla/show_bug.cgi?id=32653

[Path]
PathModified=/usr/local/bin/warsaw/core
Unit=glibc2-41-fix.service

[Install]
WantedBy=default.target
